# debian/user/mtvdownload.service
[Unit]
Description=Download der Medithekview Filmliste
After=mtvupdate.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mtvclient -D

